# debian/linode-longview.prerm
#!/bin/sh
set -e
invoke-rc.d longview stop || true
